# commandbot.service
[Unit]
Description=Command Bot

[Service]
ExecStart=/usr/bin/commandbot
Restart=on-failure

[Install]
WantedBy=multi-user.target
